# foff-milter.rc
#!/bin/sh

# PROVIDE: foff_milter
# REQUIRE: LOGIN
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable foff-milter:
# foff_milter_enable="YES"
# foff_milter_config="/usr/local/etc/foff-milter.toml"

. /etc/rc.subr

name="foff_milter"
rcvar="foff_milter_enable"

load_rc_config $name

: ${foff_milter_enable:="NO"}
: ${foff_milter_config:="/usr/local/etc/foff-milter.toml"}
: ${foff_milter_user:="foff-milter"}
: ${foff_milter_group:="foff-milter"}
: ${foff_milter_pidfile:="/var/run/foff-milter/foff-milter.pid"}

command="/usr/local/bin/foff-milter"
command_args="-c ${foff_milter_config}"
pidfile="${foff_milter_pidfile}"
required_files="${foff_milter_config}"

# FreeBSD rc.d handles backgrounding - no daemon flag needed
foff_milter_user="${foff_milter_user}"
foff_milter_group="${foff_milter_group}"

start_precmd="foff_milter_prestart"
reload_cmd="foff_milter_reload"
extra_commands="reload"

foff_milter_prestart()
{
    if [ ! -d "/var/run/foff-milter" ]; then
        install -d -o ${foff_milter_user} -g ${foff_milter_group} -m 755 /var/run/foff-milter
    fi
    if [ ! -d "/var/lib/foff-milter" ]; then
        install -d -o ${foff_milter_user} -g ${foff_milter_group} -m 755 /var/lib/foff-milter
    fi
}

foff_milter_reload()
{
    if [ -f "${pidfile}" ]; then
        pid=$(cat ${pidfile})
        if kill -0 ${pid} 2>/dev/null; then
            echo "Reloading ${name} configuration and modules..."
            if kill -HUP ${pid}; then
                echo "Configuration reload signal sent to ${name} (PID: ${pid})."
                echo "Check logs to verify successful reload."
            else
                echo "Failed to send reload signal to ${name}."
                return 1
            fi
        else
            echo "${name} PID file exists but process is not running."
            return 1
        fi
    else
        echo "${name} is not running (no PID file found)."
        return 1
    fi
}

run_rc_command "$1"
