# foff-milter.service
[Unit]
Description=FOFF Milter - Enterprise Email Security Platform
After=network.target
Wants=network.target

[Service]
Type=simple
User=foff-milter
Group=foff-milter
ExecStart=/usr/local/bin/foff-milter -c /etc/foff-milter.toml
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=foff-milter

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/foff-milter /var/run/foff-milter

[Install]
WantedBy=multi-user.target
